# service/smbd/run
#!/bin/sh -e

/usr/sbin/smbd -F -S --no-process-group
